# deploy/systemd/ibkr-strategy-runner-leaps.service
[Unit]
Description=ibkr-strategy-runner LEAPS overlay trading daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/path/to/ibkr_strategy_runner
ExecStart=/path/to/ibkr_strategy_runner/.venv/bin/python -m ibkr_strategy_runner --env-file /path/to/ibkr_strategy_runner/.env.paper-trade run-leaps --config /path/to/ibkr_strategy_runner/examples/leaps-overlay.example.json --state-dir %h/.local/state/ibkr-strategy-runner --state-backend sqlite
Restart=always
RestartSec=60
KillSignal=SIGINT

[Install]
WantedBy=default.target
